# config/plex.conf
[program:plex]
command=/usr/lib/plexmediaserver/Plex\ Media\ Server
user=plex
autostart=true
autorestart=true
stderr_logfile=/var/log/plex.err.log
stdout_logfile=/var/log/plex.out.log